# infra/systemd/lics.service
[Unit]
Description=Laboratorio de Inyección Castro Solís
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/lics
ExecStart=/opt/lics/scripts/start.sh
ExecStop=/opt/lics/scripts/stop.sh
TimeoutStartSec=300
TimeoutStopSec=120
Restart=no

[Install]
WantedBy=multi-user.target
